# debian/fglrx.postinst.in
#!/bin/sh
# Copyright (C) 2007 Mario Limonciello
# Copyright (C) 2009-2011 Canonical Ltd.

set -e

PACKAGE_NAME=#DRIVERNAME#
CVERSION=`dpkg-query -W -f='${Version}' $PACKAGE_NAME | awk -F "-" '{print $1}' | cut -d\: -f2`

ARCH=`dpkg --print-architecture`

if [ "$1" = "configure" ]; then

    #check whether libglx.so got installed.  some releases earlier than x740
    #will not actually have this libglx.so
    if [ ! -f /usr/lib/xorg/modules/extensions/libglx.so ]; then
        dpkg-divert --remove --rename --package xorg-driver-fglrx --divert /usr/lib/fglrx/libglx.so.xlibmesa /usr/lib/xorg/modules/extensions/libglx.so > /dev/null
    fi

    # Clean up any previous non-multi-arch alternatives
    if [ -n "$(update-alternatives --list gl_conf 2>/dev/null)" ]; then
        set -a $(update-alternatives --list gl_conf 2>/dev/null)
        while [ $# -ge 1 ] && [ "$#" != "configure" ]; do
            ALTERNATIVE=${1}
            update-alternatives --remove gl_conf $ALTERNATIVE
            shift
        done
    fi

    update-alternatives --force \
        --install /#SYSCONFDIR#/ld.so.conf.d/#DEB_HOST_MULTIARCH#_GL.conf #DEB_HOST_MULTIARCH#_gl_conf /#LDSOCONF# #ALTPRIORITY# \
        --slave /#DATADIR#/applications/ubuntu-amdcccle.desktop #DEB_HOST_MULTIARCH#_amdcccle_desktop /#PKGDATADIR#/amdcccle.desktop \
        --slave /#DATADIR#/applications/ubuntu-amdccclesu.desktop #DEB_HOST_MULTIARCH#_amdccclesu_desktop /#PKGDATADIR#/amdccclesu.desktop \
        --slave /#LIBDIR#/xorg/modules/drivers/fglrx_drv.so #DEB_HOST_MULTIARCH#_fglrx_drv /#PKGLIBDIR#/xorg/modules/drivers/fglrx_drv.so \
        --slave /#SYSCONFDIR#/X11/Xsession.d/10fglrx #DEB_HOST_MULTIARCH#_10fglrx /#PKGLIBDIR#/10fglrx \
        --slave /#XORGEXTRA# #DEB_HOST_MULTIARCH#_xorg_extra_modules /#PKGLIBDIR#/xorg \

    # This is for switchable graphics
    # Note: the radeon kernel module is still blacklisted
    update-alternatives --force \
        --install /#SYSCONFDIR#/ld.so.conf.d/#DEB_HOST_MULTIARCH#_GL.conf #DEB_HOST_MULTIARCH#_gl_conf /#PXLDSOCONF# #PXALTPRIORITY# \
        --slave /#DATADIR#/applications/ubuntu-amdcccle.desktop #DEB_HOST_MULTIARCH#_amdcccle_desktop /#PKGDATADIR#/amdcccle.desktop \
        --slave /#DATADIR#/applications/ubuntu-amdccclesu.desktop #DEB_HOST_MULTIARCH#_amdccclesu_desktop /#PKGDATADIR#/amdccclesu.desktop \
        --slave /#LIBDIR#/xorg/modules/drivers/fglrx_drv.so #DEB_HOST_MULTIARCH#_fglrx_drv /#PKGLIBDIR#/xorg/modules/drivers/fglrx_drv.so \
        --slave /#SYSCONFDIR#/X11/Xsession.d/10fglrx #DEB_HOST_MULTIARCH#_10fglrx /#PKGLIBDIR#/10fglrx \
        --slave /#XORGEXTRA# #DEB_HOST_MULTIARCH#_xorg_extra_modules /#PXDIR#/xorg \

    # Install just the ld.so.conf.d alternatives for the secondary architecture.
    # We need only to ensure that any mesa ld.so.conf.d alternative isn't active.
    update-alternatives --force \
        --install /#SYSCONFDIR#/ld.so.conf.d/#OTHER_ARCH#_GL.conf #OTHER_ARCH#_gl_conf /#ALTLDSOCONF# #ALTPRIORITY#
    update-alternatives --force \
        --install /#SYSCONFDIR#/ld.so.conf.d/#OTHER_ARCH#_GL.conf #OTHER_ARCH#_gl_conf /#ALTPXLDSOCONF# #PXALTPRIORITY#

    # ldconfig needs to be run immediately as we're changing /etc/ld.so.conf.d/ with
    # alternatives.
    LDCONFIG_NOTRIGGER=y ldconfig
fi

# NOTE: Usually generated by dh_installinit, but because two packages are
# providing the same initscript, this needs some manual handling.
if [ -x "/etc/init.d/atieventsd" ]; then
    update-rc.d atieventsd defaults 31 >/dev/null || exit $?
fi

#DEBHELPER#
